# daemon/weather.service
[Unit]
Description=Raspberry Pi Weather Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/weather
EnvironmentFile=/opt/weather/weather.env
ExecStart=/usr/bin/python3 -m main.weather-service \
  --db ${DB_PATH} \
  --retention ${RETENTION} \
  --bus ${BUS} \
  --bme-addr ${BME_ADDR} \
  --veml-addr ${VEML_ADDR} \
  --veml-gain ${VEML_GAIN} \
  --veml-integration-ms ${VEML_INTEGRATION_TIME} \
  --sgp-addr ${SGP_ADDR} \
  --lcd-addr ${LCD_ADDR} \
  --sample-interval ${SAMPLE_INTERVAL} \
  --display-interval ${DISPLAY_INTERVAL} \
  --host ${HOST} \
  --port ${PORT}

Restart=on-failure
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/weather /var/lib/weather
RuntimeDirectory=weather
StateDirectory=weather

[Install]
WantedBy=multi-user.target
